# bot.service
[Unit]
Description=FRKN Telegram Bot (Rust)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/frkn-bot
Environment=BOT_TOKEN=ТВОЙ_ТОКЕН
ExecStart=/root/frkn-bot/target/release/frkn-bot
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
